# cleandloads.timer
[Unit]
Description= Clean up my Downloadsdirectory

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target