# cleandloads.service
[Unit]
Description=Cleanup Downloads Directory

[Service]
Type=oneshot
ExecStart=

[Install]
WantedBy=